# build/crawler/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

*/30 * * * * source /etc/profile && /ojah/manage.py crawl
