# build/app/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

0 * * * * source /etc/profile && /ojah/manage.py stats_calculate
